# install_deps/msr_server.service
[Unit]
Description="generates data from msr in file /msr_data.toml"

[Service]
ExecStart = /usr/bin/msr_gen -r

[Install]
WantedBy=multi-user.target